# docker-deps/supervisord.ini
[supervisord]
nodaemon=true

[program:nginx]
command=/usr/sbin/nginx -c /etc/nginx/nginx.conf
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:gunicorn]
directory=/var/www
command= /usr/bin/gunicorn --name teacher_directory --bind 0.0.0.0:8000 teacher_directory.wsgi:application  --error-logfile /dev/stderr --access-logfile /dev/stdout --pid /var/run/gunicorn.pid
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0